# config/crontab
30 2 * * * ~/SetStory/bin/update-social-data
